# system/services/backup.mount
[Unit]
Description=Backup drive

[Mount]
What=LABEL=backup
Where=/media/backups
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
